# tempctl.service
[Unit]
Description=Thermal protection

[Service]
Type=idle
ExecStart=/usr/bin/tempctl

[Install]
WantedBy = multi-user.target
